# systemd/tunnel-health.timer
[Unit]
Description=Periodic Tunnel Health Check

[Timer]
OnBootSec=30
OnUnitActiveSec=60
AccuracySec=5

[Install]
WantedBy=timers.target
